# home/Library/LaunchAgents/com.dotfiles-private.sync-runtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles-private.sync-runtime</string>

    <!--
      Every 5 min: rsync the runtime subset (zsh history, REPL histories,
      ~/.claude/projects, ~/.claude/history.jsonl, …) between local and
      <gdrive>/dotfiles-runtime/${device_id}/. Pull-then-push by default.
      Each machine has its OWN ${device_id}/ subdir, so there are no
      conflicts to resolve.

      Bails silently when GDrive isn't mounted or device id is missing.
    -->
    <!--
      Redirect output inside the shell so $HOME expands per-user (see
      check-repo plist for rationale).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>mkdir -p "$HOME/Library/Logs" &amp;&amp; cd "$HOME/.dotfiles" &amp;&amp; uv run sync-private-runtime &gt;&gt; "$HOME/Library/Logs/dotfiles-sync-runtime.log" 2&gt;&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
