# home/Library/LaunchAgents/com.dotfiles-private.check-repo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles-private.check-repo</string>

    <!--
      Once an hour: git fetch + count behind/ahead/dirty in
      ~/.dotfiles-private. If any are non-zero, fire one
      macOS notification. Cooldown is 1h per distinct message; clean state
      clears the cooldown so the next problem fires immediately.

      Bails silently when the local clone is missing — safe to ship even
      on machines that haven't been bootstrapped yet.
    -->
    <!--
      Redirect both the bash shell output AND any errors into the same
      cache file the script itself writes to (~/.cache/dotfiles-private/
      check-repo.log). One file, one source of truth for "what did this
      agent do last tick?". Bash output here is normally empty — the
      script logs its own structured entries — but if the shell pipeline
      ever breaks (uv missing, $HOME unwritable, etc.) the error lands
      next to the structured logs rather than disappearing.

      Per-user $HOME expansion keeps multi-user macs (this one hosts both
      `claw` and `josh`) from interleaving log files.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>mkdir -p "$HOME/.cache/dotfiles-private" &amp;&amp; cd "$HOME/.dotfiles" &amp;&amp; uv run check-private-repo &gt;&gt; "$HOME/.cache/dotfiles-private/check-repo.log" 2&gt;&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
