# home/Library/LaunchAgents/com.dotfiles-private.detach-cloud-cache.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles-private.detach-cloud-cache</string>

    <!--
      Run detach-cloud-cache hourly. Idempotent: a run is a no-op once
      everything is already detached.

      Uses bash -lc so we get the user's PATH (mise, uv, etc.) without
      hardcoding /opt/homebrew. WorkingDirectory lets uv resolve the project
      virtualenv automatically.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME/.dotfiles" &amp;&amp; uv run detach-cloud-cache</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/dotfiles-detach-cloud-cache.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dotfiles-detach-cloud-cache.err.log</string>
</dict>
</plist>
